# src/etc/c/rc.d/pf
#!/bin/sh
#
#

daemon="/sbin/pfctl"

. /etc/rc.subr

rc_cmd $1
